# systemd/trading-bot-api.service
[Unit]
Description=Trading Bot API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=carsonodell
WorkingDirectory=/home/carsonodell/trading-bot
ExecStart=/usr/bin/python3 api/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Load environment variables from .env file
EnvironmentFile=/home/carsonodell/trading-bot/.env

[Install]
WantedBy=multi-user.target
